# microproyecto_balanceo/systemd/web-3000.service
[Unit]
Description=Web App 3000
After=network.target

[Service]
Environment=PORT=3000
ExecStart=/usr/bin/node /opt/app/server.js
Restart=always

[Install]
WantedBy=multi-user.target